# certbot.timer
[Unit]
Description=Three month renewal of Let's Encrypt's certificates

[Timer]
OnCalendar= 2m 3w
Persistent=true

[Install]
WantedBy=timers.target
